# rust/aggregator/config/crontab
*/10 * * * * redis-cli -h host.docker.internal LPUSH aggregator:worker:jobs run:all
